# home/.config/autostart/discord.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Discord
Exec=discord
Terminal=false
StartupNotify=false
