# home/.config/autostart/emacs.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Emacs
Exec=emacs
Terminal=false
StartupNotify=true
