# home/.config/autostart/git-dude.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Git-dude
Exec=git-dude ~/.git-dude
Terminal=false
StartupNotify=false
