# home/.config/autostart/librewolf.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Librewolf
Exec=librewolf
Terminal=false
StartupNotify=false
