# home/.config/autostart/raindrop.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Raindrop
Exec=raindrop
Terminal=false
StartupNotify=false
