# home/.config/autostart/syngestures.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Syngestures
Exec=syngestures
Terminal=false
StartupNotify=false
