# home/.config/autostart/unclutter.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Unclutter
Exec=unclutter
Terminal=false
StartupNotify=false
