# home/.config/autostart/urxvtcd.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Urxvtcd
Exec=urxvtcd
Terminal=false
StartupNotify=false
